# ansible/playbooks/templates/bliss-server.service.j2
[Unit]
Description=Bliss Demo Servers ({{ vm_context }})
After=network.target
Before=bliss-session-{{ vm_context }}.service

[Service]
Type=simple
User={{ service_user }}
{% if service_group %}
Group={{ service_group }}
{% endif %}
WorkingDirectory={{ install_base_path }}/bliss

# Start command - activate conda environment and run BLISS servers
ExecStart=/bin/bash -c 'HOSTNAME=$(hostname) && export BEACON_HOST=$HOSTNAME:{{ bliss_beacon_port }} && export TANGO_HOST=$HOSTNAME:{{ bliss_tango_port }} && source {{ install_base_path }}/miniconda/bin/activate bliss_dev && {{ bliss_server_command }}'

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bliss-server-{{ vm_context }}

[Install]
WantedBy=multi-user.target
