# ansible/playbooks/templates/mxcubeweb.service.j2
[Unit]
Description=MXCubeWeb ({{ vm_context }})
After=network.target

[Service]
Type=simple
User={{ service_user }}
{% if service_group %}
Group={{ service_group }}
{% endif %}
WorkingDirectory={{ install_base_path }}/mxcubeweb

# Environment variables
Environment="PYTHONPATH={{ install_base_path }}/mxcubecore:{{ install_base_path }}/mxcubeweb"
Environment="MXCUBE_HOST={{ mxcubeweb_config.host }}"
Environment="MXCUBE_PORT={{ mxcubeweb_config.port }}"

# Start command - activate conda environment and run mxcubeweb-server
ExecStart=/bin/bash -c 'source {{ install_base_path }}/miniconda/bin/activate mxcubeweb && mxcubeweb-server --ra -r {{ install_base_path }}/mxcubeweb/{{ mxcube_config_dir }}/ --static-folder {{ install_base_path }}/mxcubeweb/ui/build/ -L {{ mxcube_log_level }} -l {{ log_base_path }}/mxcube.log'

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mxcubeweb-{{ vm_context }}

[Install]
WantedBy=multi-user.target
